# monitor_temp
### BEGIN INIT INFO 
# Provides: automount 
# Required-Start: $remote_fs $syslog 
# Required-Stop: $remote_fs $syslog 
# Default-Start: 2 3 4 5 
# Default-Stop: 0 1 6 
# Short-Description: Start daemon at boot time 
# Description: Enable service provided by daemon. 
### END INIT INFO 
#! /bin/sh
# /etc/init.d/blah
#

# Some things that run always

# Carry out specific functions when asked to by the system
case "$1" in
  start)
     echo "Started temperature logging"
     python -u /var/temperature/monitor.py > /dev/null&
    ;;
  stop)
      echo "Stopped temperature logging"
      pkill -f "python -u /var/temperature/monitor.py"
    ;;
  log)
      echo "Start log of temperatures to database and file: /var/temperature/temps.log"
      python -u /var/temperature/monitor.py >>/var/temperature/temps.log&
    ;;
  *)
    echo "Usage: service monitor_temp [start stop log]"
    exit 1
    ;;
esac

exit 0
